# deployments/systemd/ai-stack-gateway.service
[Unit]
Description=AI Stack API Gateway Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ai-stack
Group=ai-stack
WorkingDirectory=/opt/ai-stack/gateway
Environment="PATH=/opt/ai-stack/gateway/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PORT=9000"
Environment="HOST=0.0.0.0"

# 启动命令
ExecStart=/opt/ai-stack/gateway/venv/bin/uvicorn main:app --host 0.0.0.0 --port 9000 --workers 4

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-stack-gateway

[Install]
WantedBy=multi-user.target

